# auto-cleaner.service
[Unit]
Description=Auto Cleaner Cache Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-cleaner.sh
